# local/etc/cron.d/kbnet.com-zfs-report.orig
MAILTO="itgroup@mail-szerver.localdomain"
MAILFROM="cron@zfs-szerver.localdomain"

0 23 * * *  root  /opt/kbnet.com-zfs-report/usr/share/kbnet.com-zfs-report.bash
